# scripts/systemd/octoagent-local.service
[Unit]
Description=OctoAgent Local Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=sieve-pub
WorkingDirectory=/home/sieve-pub/public-workspace/octoagent
Environment=PATH=/snap/bin:/usr/bin:/usr/sbin:/bin:/home/sieve-pub/public-workspace/octoagent/scripts
Environment=TMPDIR=/home/sieve-pub/public-workspace/octoagent/tmp
Environment=OCTOAGENT_PYTHON_BIN=/home/sieve-pub/public-workspace/octoagent/backend/.venv/bin/python
Environment=OCTO_AGENT_CAPABILITY_SOURCE=/home/sieve-pub/public-workspace/octoagent
Environment=OCTOAGENT_MANAGED_BY_SYSTEMD=1
Environment=BG_JOB_ISOLATED_LOOPS=true
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/sieve-pub/public-workspace/octoagent/scripts/start-octoagent.sh run
ExecStop=/home/sieve-pub/public-workspace/octoagent/scripts/start-octoagent.sh stop-runtime
Restart=on-failure
RestartSec=10
KillMode=mixed
TimeoutStartSec=900
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
