# deploy/octoagent-execution-worker.service
[Unit]
Description=OctoAgent Independent Execution Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sieve-pub
Group=sieve-pub
WorkingDirectory=/home/sieve-pub/public-workspace/octoagent/backend
EnvironmentFile=-/etc/octoagent/execution-worker.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/sieve-pub/public-workspace/octoagent/backend
ExecStart=/home/sieve-pub/public-workspace/octoagent/backend/.venv/bin/python /home/sieve-pub/public-workspace/octoagent/backend/scripts/run_execution_worker.py --host ${OCTO_EXECUTION_WORKER_HOST} --port ${OCTO_EXECUTION_WORKER_PORT} --node-id ${OCTO_EXECUTION_WORKER_NODE_ID} --token ${OCTO_EXECUTION_WORKER_TOKEN} --capacity ${OCTO_EXECUTION_WORKER_CAPACITY} --callback-retries ${OCTO_EXECUTION_WORKER_CALLBACK_RETRIES}
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/sieve-pub/public-workspace/octoagent/workspace /home/sieve-pub/public-workspace/octoagent/logs /home/sieve-pub/public-workspace/octoagent/tmp

[Install]
WantedBy=multi-user.target
